# source/vsm/etc/systemd/vsm.target
[Unit]
Description=vsm target allowing to start/stop all vsm instances at once
[Install]
WantedBy=multi-user.target
